# lessons/146/ansible/roles/myapp/files/myapp.service
[Unit]
Description=My App
Wants=network-online.target
After=network-online.target

StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
User=myapp
Group=myapp
Type=simple
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/web

[Install]
WantedBy=multi-user.target
